# library/centos/roles/prometheus-node-exporter/templates/node_exporter.upstart.j2
description "Prometheus node exporter"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn
respawn limit 10 5
setuid {{ prometheus_n_e_user }}
setgid {{ prometheus_n_e_user }}

script
  exec {{ prometheus_n_e_cmd }} {{ prometheus_n_e_opts }} {{ prometheus_n_e_additional_opts }} > {{ prometheus_n_e_logdir }}/node_exporter.log 2>&1
end script
